# openports/quassel/quassel/pkg/quasselcore.rc
#!/bin/ksh

daemon="${TRUEPREFIX}/bin/quasselcore"
daemon_flags="-c ${VARBASE}/quassel"
daemon_user=_quassel

. /etc/rc.d/rc.subr

rc_cmd $1
